# Hyprland-Config/.config/systemd/user/mybar.service
[Unit]
Description=AGS Status Bar
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/ags run %h/.config/ags/status-bar/app.ts
Restart=on-failure

[Install]
WantedBy=graphical-session.target
